# ops/systemd/bountyverdict-directory-monitor.timer
[Unit]
Description=Check BountyVerdict agent directories hourly

[Timer]
OnActiveSec=15min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=true
Unit=bountyverdict-directory-monitor.service

[Install]
WantedBy=timers.target
